# crontab
45 14 * * * /usr/local/bin/vim -s /home/admin/dotfiles/update4vimplugin
